# root/etc/services.d/raneto/run
#!/usr/bin/with-contenv bash

cd /app/raneto || exit

DEBUG=raneto exec \
	s6-setuidgid abc node server.js
