# root/etc/s6-overlay/s6-rc.d/init-raneto-config/up
/etc/s6-overlay/s6-rc.d/init-raneto-config/run
